# deploy/systemd/deploy-hook.service
[Unit]
Description=Transcript Library — Deploy Hook Listener
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/transcript-library/current
ExecStart=/usr/bin/node --import tsx deploy/deploy-hook.ts
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

# Read the shared env file for DEPLOY_WEBHOOK_SECRET
EnvironmentFile=/srv/transcript-library/.env.local

StandardOutput=journal
StandardError=journal
SyslogIdentifier=deploy-hook

[Install]
WantedBy=multi-user.target
